# docker/supervisor.conf.d/mysql.conf
[program:mysql]
autostart=true
autorestart=true
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/sbin/mysqld
